# archive/com.jonsims.dailynewsagent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.jonsims.dailynewsagent</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/jonsims/Documents/Obsidian Vault/Daily Updates (Agent)/daily-news-agent/run_daily_update.sh</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>4</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/jonsims/Documents/Obsidian Vault/Daily Updates (Agent)/daily-news-agent/logs/stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/jonsims/Documents/Obsidian Vault/Daily Updates (Agent)/daily-news-agent/logs/stderr.log</string>

    <key>WorkingDirectory</key>
    <string>/Users/jonsims/Documents/Obsidian Vault/Daily Updates (Agent)/daily-news-agent</string>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>